# gcode-monitor.service
# GCode Monitor Systemd Service Template
# This file contains placeholders that are replaced during installation.
# Run install_service.sh to generate and install the actual service file.
#
# Placeholders:
#   %USERNAME% - Your system username
#   %HOME_DIR% - Your home directory
#   %PROJECT_DIR% - Pi-gcode-server installation directory
#   %WATCH_DIR% - Directory to monitor for .gcode files
#   %LOCAL_SUBNET% - Auto-detected local network subnet (e.g., 192.168.1.0/24)

[Unit]
Description=GCode File Monitor and Sync Service
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=600

[Service]
Type=simple
User=%USERNAME%
WorkingDirectory=%PROJECT_DIR%
ExecStart=/usr/bin/python3 %PROJECT_DIR%/monitor_and_sync.py

# Restart policy
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax=500M
CPUQuota=50%
TasksMax=20

# Network hardening (restrict to local subnet only)
RestrictAddressFamilies=AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=%LOCAL_SUBNET%

# Filesystem protection (minimal access)
NoNewPrivileges=true
# PrivateTmp=true  # Disabled: can cause namespace issues
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%HOME_DIR%/.gcode_sync.log
ReadWritePaths=%HOME_DIR%/.local
ReadOnlyPaths=%WATCH_DIR%
ReadOnlyPaths=%PROJECT_DIR%

# Additional security hardening
# PrivateDevices=true  # Disabled: can cause namespace issues
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# RestrictNamespaces=true  # Disabled: incompatible with Python applications
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectKernelLogs=true
ProtectClock=true

# Capability restrictions
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering - relaxed for Python compatibility
# SystemCallFilter=@system-service
# SystemCallFilter=~@privileged @resources @obsolete
# SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gcode-monitor

[Install]
WantedBy=multi-user.target
